# rpi-sensor-block/services/lz.ard/lz.ard.service
[Unit]
Description=lz's main ard reader.
After=network.target auditd.service lz.multicast.service

[Service]
ExecStart=/opt/lz/rpi-sensor-block/services/lz.ard/ard.pl
# ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=lz.ard
